# Client/Extra/heatcontroller.service
[Unit]
Description=Heat Controller

[Service]
Type=simple
WorkingDirectory="<WORKPATH>"
ExecStart=<EXECSTART>

[Install]
WantedBy=multi-user.target