# Client/Extra/heatcontroller.service
[Unit]
Description=Heat Controller

[Service]
Type=simple
RootDirectory=""
ExecStart=<INSERT PATH>

[Install]
WantedBy=multi-user.target